# supervisor_celery_monitor.example.conf
; Supervisor config for celery_monitor daemon
; Copy to /etc/supervisor/conf.d/ or include in supervisord.conf

[program:celery_monitor]
command=/path/to/venv/bin/celery_monitor -A your_project.celery:app --config /path/to/config.yaml --daemon
directory=/path/to/your_project
user=celery
environment=CELERY_APP="your_project.celery:app"
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/celery_monitor.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
